# deploy/agents-cron.timer
[Unit]
Description=Agents Job Cron Timer (every minute)

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target
