# deploy/agents-ws.service
[Unit]
Description=Agents WebSocket Server (Instagram Browser Sessions)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/agents/portal
ExecStart=/usr/bin/doppler run -p agents -c prd -- /usr/bin/npx tsx server/ws-server.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
